# examples/systemd/osm-notes-ingestion-daemon.service
[Unit]
Description=OSM Notes Ingestion Processing Daemon
Documentation=file:///home/notes/OSM-Notes-Ingestion/docs/Process_API.md
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
# Production: Run as 'notes' user (different from login user 'angoca')
# Verify user exists: getent passwd notes
# Find primary group: id notes | grep -o 'gid=[0-9]*([^)]*)' | cut -d'(' -f2 | cut -d')' -f1
User=notes
# Group is optional - systemd will use user's primary group automatically
# If you need to specify, find it with: id -gn notes
# Group=notes
WorkingDirectory=/home/notes/OSM-Notes-Ingestion

# Environment variables
# Ensure PATH includes standard directories for commands (setsid, readlink, etc.)
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LOG_LEVEL=INFO"
Environment="DAEMON_SLEEP_INTERVAL=60"
Environment="CLEAN=true"
Environment="SKIP_XML_VALIDATION=true"
Environment="SKIP_CSV_VALIDATION=true"

# Executable
# Use absolute path with bash explicitly
ExecStart=/bin/bash /home/notes/OSM-Notes-Ingestion/bin/process/processAPINotesDaemon.sh
ExecReload=/bin/kill -HUP $MAINPID

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

# Resource limits (optional)
# LimitNOFILE=65536
# MemoryMax=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=osm-notes-ingestion-daemon

[Install]
WantedBy=multi-user.target
